# setup/com.flyfx.daily-scan.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.flyfx.daily-scan</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>REPLACE_WITH_YOUR_HOME/Desktop/FlyFX/deals-machine/daily-scan.sh</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>5</integer>
        <key>Minute</key>
        <integer>30</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>REPLACE_WITH_YOUR_HOME/Desktop/FlyFX/deals-machine/logs/launchd-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>REPLACE_WITH_YOUR_HOME/Desktop/FlyFX/deals-machine/logs/launchd-stderr.log</string>

    <key>WorkingDirectory</key>
    <string>REPLACE_WITH_YOUR_HOME/Desktop/FlyFX/deals-machine</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>

    <!-- Run on load if the scheduled time was missed (Mac was asleep) -->
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
